# lock.target
[Unit]
Conflicts=unlock.target
Description=Lock the current session
