# sleep.target
[Unit]
Description=User-level target triggered when the system is about to sleep.
