# unlock.target
[Unit]
Conflicts=lock.target
Description=Unlock the current session
